# systemd/gpisend.service.in
[Unit]
Description=GPI network messenger
After=network.target remote-fs.target nss-lookup.target
ConditionPathExists=/etc/gpisend.conf

[Service]
Type=simple
ExecStart=@prefix@/bin/gpisend
PrivateTmp=true
Restart=always
RestartSec=2
StartLimitInterval=120
StartLimitBurst=50

[Install]
WantedBy=multi-user.target
