# debian/postinst
#!/bin/sh
# postinst script for gpisend

case "$1" in
    configure)
      ldconfig
      if test ! -f /etc/gpisend.conf ; then
	  cp /usr/share/gpisend/gpisend-sample.conf /etc/gpisend.conf
      fi
      /bin/systemctl daemon-reload
      /bin/systemctl enable gpisend.service
      /bin/systemctl restart gpisend.service
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
